# root/etc/s6-overlay/s6-rc.d/svc-warp-monitor/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "**** Starting WARP connection monitor ****"

WARP_CONF="/config/wgcf0.conf"

while true; do
    sleep 10

    # Check if WARP interface is up
    if ! ip link show wgcf0 2>/dev/null | grep -q "<.*UP.*>"; then
        echo "**** CRITICAL: WARP connection lost! ****"
        echo "**** Kill switch is active - no traffic can leak ****"
        echo "**** Attempting to reconnect... ****"

        # Try to bring interface back up (manual setup to avoid sysctl)
        ip link del wgcf0 2>/dev/null || true
        ip link add dev wgcf0 type wireguard
        PRIVATE_KEY=$(grep "^PrivateKey" "${WARP_CONF}" | cut -d= -f2- | tr -d ' ')
        wg set wgcf0 private-key <(echo "$PRIVATE_KEY")
        PEER_KEY=$(grep "^PublicKey" "${WARP_CONF}" | cut -d= -f2- | tr -d ' ')
        ENDPOINT=$(grep "^Endpoint" "${WARP_CONF}" | cut -d= -f2- | tr -d ' ')
        ALLOWED_IPS=$(grep "^AllowedIPs" "${WARP_CONF}" | cut -d= -f2- | tr -d ' ')
        wg set wgcf0 peer "$PEER_KEY" endpoint "$ENDPOINT" allowed-ips "$ALLOWED_IPS" persistent-keepalive 25
        IPV4_ADDR=$(grep "^Address" "${WARP_CONF}" | awk -F'[, ]' '{print $3}')
        IPV6_ADDR=$(grep "^Address" "${WARP_CONF}" | awk -F'[, ]' '{print $4}')
        [ -n "$IPV4_ADDR" ] && ip -4 address add "$IPV4_ADDR" dev wgcf0
        [ -n "$IPV6_ADDR" ] && ip -6 address add "$IPV6_ADDR" dev wgcf0
        MTU=$(grep "^MTU" "${WARP_CONF}" | awk '{print $3}')
        [ -n "$MTU" ] && ip link set mtu "$MTU" dev wgcf0
        ip link set wgcf0 up
        ip route del default 2>/dev/null || true
        ip route add default dev wgcf0
        sleep 5

        # If still not connected after retry, this is critical
        if ! ip link show wgcf0 2>/dev/null | grep -q "<.*UP.*>"; then
            echo "**** FATAL: WARP reconnection failed ****"
            echo "**** Container will shutdown to prevent traffic leaks ****"
            # Send termination signal to s6-overlay
            s6-svscanctl -t /var/run/s6/services
        else
            echo "**** WARP reconnected successfully ****"
        fi
    fi
done
